# week3/myCoolTarget.target
[Unit]
Description=cool target
Requires=bionic.service noble.service jammy.service
